# CPSC-334-Linux&DevOps/systemd-counter-JustFady/Counter.service
[Unit]
Description=Simple Counter Service
After=network.target

[Service]
ExecStart=/usr/bin/python3 /usr/local/bin/counter.py
Restart=always
User=counter
Group=counter

[Install]
WantedBy=multi-user.target
